# w767-os/audit/gbs-audit-quick.timer
[Unit]
Description=Run gbs-audit-quick every 5 minutes
After=multi-user.target

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=15s
Unit=gbs-audit-quick.service
Persistent=false

[Install]
WantedBy=timers.target
