# w767-os/audit/gbs-audit-quick.service
[Unit]
Description=Galaxy Book S quick hot-files audit
ConditionPathExists=/usr/local/bin/gbs-audit-quick

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gbs-audit-quick
Nice=15
IOSchedulingClass=idle
